# supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
command=service mysql start

[program:node]
command=node /usr/src/question-app/api/index.js